# scripts/clash-utun.timer
[Unit]
Description=Run clash-utun on boot

[Timer]
OnBootSec=30s

[Install]
WantedBy=timers.target